# pop/systemd/restic-maintain.timer
[Unit]
Description=Restic maintenance

[Timer]
OnUnitActiveSec=5m
OnBootSec=5m

[Install]
WantedBy=timers.target
